# deploy/systemd/mcchallonge.service
[Unit]
Description=McChallonge Flask Dashboard (Waitress)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/mcchallonge
EnvironmentFile=/opt/mcchallonge/.env
ExecStart=/opt/mcchallonge/venv/bin/python -m mcchallonge.app serve --port 8000 --threads 8
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
